# Sources/ospboard/opt/osp/src/wcn36xx-monitor/wcn36xx-monitor.service
[Unit]
Description=wcn36xx firmware bug monitor

[Service]
ExecStart= wcn36xx-monitor
Restart=always

[Install]
WantedBy=multi-user.target
